# deploy/systemd/alloc-context-ingest.service
[Unit]
Description=AllocContext ingest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=trading
Group=trading
WorkingDirectory=/opt/alloc-context
EnvironmentFile=-/opt/alloc-context/.env
ExecStart=/opt/alloc-context/.venv/bin/python -m alloccontext ingest
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
