# deploy/systemd/alloc-context-alerts.timer
[Unit]
Description=AllocContext hourly band alerts
Requires=alloc-context-alerts.service

[Timer]
OnCalendar=*-*-* *:10:00
Persistent=true

[Install]
WantedBy=timers.target
