# deploy/systemd/alloc-context-backup.timer
[Unit]
Description=AllocContext SQLite backup timer (daily)
Requires=alloc-context-backup.service

[Timer]
OnCalendar=*-*-* 04:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
